# railway/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:api]
directory=/app/apps/api
command=/opt/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port %(ENV_INTERNAL_API_PORT)s --proxy-headers
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app/apps/api"

[program:web]
directory=/app/apps/web
command=/usr/bin/env PORT=%(ENV_INTERNAL_WEB_PORT)s HOSTNAME=127.0.0.1 node /app/apps/web/server.js
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
